# image/rootfs/etc/services.d/sshd/run
#!/usr/bin/execlineb -P
foreground { s6-mkdir -p -m 750 /var/run/sshd }
fdmove -c 2 1
if { /usr/sbin/sshd -t }
/usr/sbin/sshd -D -e
